# configs/os/services/kvmd-flash-pico@.service
[Service]
Type=oneshot
ExecStart=/usr/lib/kvmd/kvmd-udev-flash-pico %I
